# .docker/supervisord.ini
[supervisord]
nodaemon=true

[program:octane]
command= php /var/www/app/artisan serve --host=0.0.0.0 --port=80
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
